# src/upnpc-daemon.service.in
[Unit]
Description=redirecting ports to localhost by UPnP

[Service]
Type=oneshot
